# solarmon.service
[Unit]
Description=Solarmon
After=network.target
StartLimitIntervalSec=10min
StartLimitBurst=3

[Service]
# Sync Python deps before each start so a git pull that changes requirements.txt
# is applied automatically. '+' runs as root (libs live in system dist-packages
# while the service runs as User=pi); '-' tolerates failure (e.g. offline boot).
ExecStartPre=-+/usr/bin/pip3 install -r requirements.txt
ExecStart=/usr/bin/python3 -u solarmon.py
WorkingDirectory=/home/pi/solarmon
StandardOutput=inherit
StandardError=inherit
RestartSec=10
RestartPreventExitStatus=0
WatchdogSec=60
Restart=always
User=pi

[Install]
WantedBy=multi-user.target
